# ai_chat/deploy/itch-games-ai-chat-qwertystock.service
[Unit]
Description=Orbital Courier AI development chat on qwertystock.com
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/generic/itch_games
Environment=AI_CHAT_HOST=127.0.0.1
Environment=AI_CHAT_PORT=8765
Environment=AI_CHAT_DATA_DIR=/home/generic/itch_games/ai_chat/data
EnvironmentFile=-/etc/itch-games-ai-chat.env
Environment=AI_CHAT_DEPLOY_SCRIPT=/home/generic/itch_games/ai_chat/deploy/itch-games-autodeploy-qwertystock.sh
Environment=AI_CHAT_DEPLOY_RELAY_URLS=https://moscow-io-games.mecharulez.com/ai_chat/api/deploy-relay
Environment=AI_CHAT_DEPLOY_RELAY_TIMEOUT_MS=5000
Environment=UNSOCCER_AUTOSTART=1
Environment=UNSOCCER_PORT=8787
ExecStart=/usr/bin/node /home/generic/itch_games/ai_chat/server_node.js
Restart=always
RestartSec=3
User=generic
Group=generic

[Install]
WantedBy=multi-user.target
